# Citadel-01/operations/deployment/systemd-services/citadel-gateway-simple.service
[Unit]
Description=Citadel LLM API Gateway - Production with Enhanced Restart
Documentation=https://citadel.local/docs
After=network-online.target ollama.service redis-server.service
Wants=network-online.target
Requires=ollama.service redis-server.service
BindsTo=ollama.service

[Service]
Type=exec
User=agent0
Group=citadel
WorkingDirectory=/opt/citadel

# Environment Configuration
Environment="CITADEL_HOME=/opt/citadel"
Environment="CITADEL_ENV=production"
Environment="PYTHONPATH=/opt/citadel/src"
Environment="PATH=/opt/citadel/citadel_venv/bin:/usr/local/bin:/usr/bin:/bin"

# Production service with proper port configuration
ExecStart=/opt/citadel/citadel_venv/bin/uvicorn citadel_llm.api.gateway:app --host 0.0.0.0 --port 8002 --workers 8 --log-level info --access-log

# Enhanced Health check and monitoring
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStartSec=120
TimeoutStopSec=60

# Enhanced Restart Configuration
Restart=always
RestartSec=5

# Enhanced Start Limits with longer interval and more attempts
StartLimitInterval=1800  # 30 minutes
StartLimitBurst=10       # Allow up to 10 restarts in 30 minutes

# Security Configuration
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/citadel/logs /opt/citadel/var
PrivateTmp=true

# Resource Limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=8G
CPUQuota=800%

# Logging Configuration
StandardOutput=append:/opt/citadel/logs/api-gateway/service.log
StandardError=append:/opt/citadel/logs/api-gateway/error.log
SyslogIdentifier=citadel-gateway

[Install]
WantedBy=multi-user.target
Also=redis-server.service ollama.service
